# target/linux/ath79/nand/base-files/etc/hotplug.d/firmware/10-ath9k-eeprom
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"ath9k-eeprom-ahb-18100000.wmac.bin")
	case $board in
	*)
		caldata_die "board $board is not supported yet"
		;;
	esac
	;;
esac
